# rtl433-us/rootfs/etc/services.d/log/run
#!/usr/bin/with-contenv sh
exec svlogd -tt /var/log/rtl433
