# rtl433-us/rootfs/etc/services.d/rtl433/run
#!/usr/bin/with-contenv sh
exec /run.sh
